# brightness-init.service
[Unit]
Description=Enable write access to screen brightness

[Service]
Type=oneshot
ExecStart=/usr/bin/simple-brightness init

[Install]
WantedBy=default.target
